# .config/systemd/user/lxpolkit.service
[Unit]
Description=lxpolkit
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/lxpolkit
Restart=on-failure
