# .config/systemd/user/picom.service
[Unit]
Description=picom, standalone compositor for Xorg 
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/picom
Restart=on-failure
